# roles/elk/files/deconst-curator.timer
[Unit]
Description=Deconst Curator Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
