# roles/build/files/services/deconst-data-letsencrypt.service
[Unit]
Description=Deconst Let's Encrypt Certificate Data
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/opt/bin/datavolumecontainer \
  --name deconst-data-letsencrypt \
  --volume /etc/letsencrypt \
  --volume /var/lib/letsencrypt \
  --volume /cache \
  echo
